# mkosi.images/base/mkosi.extra/etc/rc.local
#!/bin/sh
set -ue

# Copy package metadata from /usr/share/factory asynchronously.
cp -rpf --update=all /usr/share/factory/var / &

# Fix services imported by systemd-confext in /etc not being loaded.
systemctl daemon-reload
# This is intended to fix some problems caused by systemd-sysext.
# When systemd-sysext merge sysext images, it won't trigger those services inside to be activated.
systemctl --no-block default
