# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/systemd-homed-firstboot.service.d/directory.conf
[Service]
ExecStart=
ExecStart=homectl firstboot --prompt-new-user --storage=directory
