# saved/desktop/mkosi.extra/usr/lib/systemd/system/rtkit-daemon.service.d/10-fix-ulimit.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/rtkit-daemon --no-limit-resources
# Let's use the ulimit interface provided by systemd.
LimitFSIZE=0
LimitMEMLOCK=0
LimitMSGQUEUE=0
LimitNICE=0
LimitNOFILE=50
# Changed from 3 to 4 because 3 sometimes makes the daemon fail to start and I really don't know why.
LimitNPROC=4
LimitRTPRIO=0
LimitRTTIME=200000
