# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/lightdm.service.d/chown.conf
[Service]
ExecStartPre=/usr/bin/mkdir -p /var/lib/lightdm
ExecStartPre=/usr/bin/chown -R lightdm:lightdm /var/lib/lightdm
